# quadlets/redis.container
# Podman Quadlet for the Redis service

[Unit]
Description=Redis container for SheepVibes
Wants=network.target
After=network.target
# Depends on the volume service being started first
Wants=sheepvibes-redis-volume-volume.service
After=sheepvibes-redis-volume-volume.service

[Container]
ContainerName=sheepvibes-redis
Image=docker.io/redis:alpine
Network=sheepvibes.network
# Mount the named volume into the container's data directory
Volume=sheepvibes-redis-volume:/data

# Expose the Redis port to the host for debugging or external tools if needed.
# For container-to-container communication, this is not strictly necessary.
PublishPort=127.0.0.1:6379:6379

[Service]
Restart=always
TimeoutStartSec=90

# Before starting, ensure the image is pulled. This prevents race conditions on a fresh system setup.
ExecStartPre=/usr/bin/podman pull docker.io/redis:alpine
