# quadlets/sheepvibes-app.container
# Podman Quadlet for the SheepVibes application

[Unit]
Description=SheepVibes Application Container

# Depend on sheepvibes-redis
Requires=sheepvibes-redis.service
After=sheepvibes-redis.service

[Container]
# Use the locally built image
Image=ghcr.io/sheepdestroyer/sheepvibes:latest
ContainerName=sheepvibes-app
Network=sheepvibes-network.network
PublishPort=127.0.0.1:5000:5000
Volume=sheepvibes-db.volume:/app/data

# Environment variables
Environment=DATABASE_PATH=/app/data/sheepvibes.db
# Connect to Redis using its container name as the hostname
Environment=CACHE_REDIS_URL=redis://sheepvibes-redis:6379/0
Environment=FLASK_APP=backend.app
Environment=PYTHONPATH=/app
Environment=UPDATE_INTERVAL_MINUTES=15
Environment=FLASK_RUN_HOST=0.0.0.0

[Service]
Restart=on-failure
TimeoutStartSec=90
